# package/deb/prerm
#!/bin/bash
set -e

SERVICE_NAME=389-ds-exporter.service

if command -v systemctl >/dev/null 2>&1; then
    echo "Stopping $SERVICE_NAME..."
    systemctl stop "$SERVICE_NAME" || true
    systemctl disable "$SERVICE_NAME" || true
fi

exit 0
